# packaging/systemd/p2p-offer.service
[Unit]
Description=WebRTC P2P Tunnel Offer Service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=p2ptunnel
Group=p2ptunnel
UMask=0077

ExecStartPre=/usr/local/bin/p2pctl check-config --config /etc/p2ptunnel/offer/config.toml
ExecStart=/usr/local/bin/p2p-offer run --config /etc/p2ptunnel/offer/config.toml

Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s
KillSignal=SIGTERM

StateDirectory=p2ptunnel-offer
LogsDirectory=p2ptunnel-offer

StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true

[Install]
WantedBy=multi-user.target
