# packaging/systemd/p2p-offer-notify.service
[Unit]
Description=WebRTC P2P Tunnel Offer Service (with sd_notify readiness)
Wants=network-online.target
After=network-online.target

[Service]
# Requires a p2p-offer binary built with --features sd-notify (see
# docs/SYSTEMD.md, "Optional: sd_notify readiness (Type=notify)"). A binary
# built without that feature never sends READY=1, so it would time out here
# under the default TimeoutStartSec — use p2p-offer.service (Type=simple)
# unless you specifically built the sd-notify-enabled binary.
Type=notify
NotifyAccess=main
User=p2ptunnel
Group=p2ptunnel
UMask=0077

ExecStartPre=/usr/local/bin/p2pctl check-config --config /etc/p2ptunnel/offer/config.toml
ExecStart=/usr/local/bin/p2p-offer run --config /etc/p2ptunnel/offer/config.toml

Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s
KillSignal=SIGTERM

StateDirectory=p2ptunnel-offer
LogsDirectory=p2ptunnel-offer

StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true

[Install]
WantedBy=multi-user.target
